# systemd/quant.service
[Unit]
Description=A股量化数据服务
Documentation=https://github.com/quant-service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/quant_service
Environment=PYTHONPATH=/root/quant_service
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /root/quant_service/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
